# spliit/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Spliit
# Runs the Nginx daemon
# ==============================================================================

# Wait for spliit to become available
bashio::net.wait_for 3000

bashio::log.info "Starting NGinx...."

whoami
ls -la /var/lib/nginx
ls -la /var/log/nginx

exec nginx